# pkg/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable ddnuts.service
systemctl start ddnuts.service
